# docker-php5full/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:php5]
command=php-fpm

[program:ssh]
command=/usr/sbin/sshd -D